# deploy/climateumbral-deploy.timer
[Unit]
Description=ClimateUmbral auto-deploy check, every 2 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
RandomizedDelaySec=20

[Install]
WantedBy=timers.target
